# src/main/resources/docker/16-2/eportal/supervisord.conf
[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
nodaemon=true
childlogdir=/var/log/

[program:setup-prerequisites.sh]
startsecs=0
command=/usr/local/nom/sbin/setup-prerequisites.sh
autorestart=unexpected
exitcodes=0,2

[program:eportal]
command=/usr/local/nom/sbin/eportal -f
autorestart=true

[program:set_ip.sh]
startsecs=0
command=/usr/local/nom/sbin/set_ip.sh

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nom-nginx]
startsecs=0
command=/usr/local/nom/sbin/nginx -c /usr/local/nom/etc/nom-nginx/nginx.conf

##SSO##[program:set_sso_ui.sh]
##SSO##startsecs=0
##SSO##command=/usr/local/nom/sbin/set_sso_ui.sh
