# src/main/resources/docker/16-2/openldap/supervisord.conf
[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
nodaemon=true
childlogdir=/var/log/

[program:openldap.sh]
startsecs=1
command=/usr/local/openldap.sh

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
