# src/main/resources/docker/16-2/shibboleth/supervisord.conf
[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
nodaemon=true
childlogdir=/var/log/

[program:shibboleth.sh]
startsecs=1
command=/usr/local/shibboleth.sh

[program:set_ip.sh]
startsecs=0
command=/usr/local/nom/sbin/set_ip.sh

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
